# renew-bot.service
[Unit]
Description=Renew Bot Telegram bot
After=network.target

[Service]
Type=simple
# User that will run the bot. Adjust if you do not want to run as root.
User=root
# Working directory where bot.py and the .env file reside
WorkingDirectory=/root/renew-bot
# Run the bot using the virtual environment's python
ExecStart=python3 bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
